# config/launchctl/com.plex.pooky.email.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>com.plex.pooky.email</string>
        <key>ProgramArguments</key>
        <array>
                <string>/usr/local/bin/plexReport.sh</string>
        </array>
        <key>StartCalendarInterval</key>
        <dict>
                <key>Hour</key>
                <integer>12</integer>
                <key>Minute</key>
                <integer>0</integer>
                <key>Weekday</key>
                <integer>5</integer>
        </dict>
</dict>
</plist>
